# packages/kmod-5.10-nvidia/nvidia-persistenced.service.in
[Unit]
Description=NVIDIA Persistence Daemon
After=load-tesla-kernel-modules.service load-open-gpu-kernel-modules.service

[Service]
Type=forking
ExecStart=__NVIDIA_BINDIR__/nvidia-persistenced --user nvidia --verbose

[Install]
RequiredBy=preconfigured.target
